# hassio-addon-jethub-mqtt-io/rootfs/etc/services.d/mqtt-io/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting mqtt-io.."

/home/mqtt_io/venv/bin/python -m mqtt_io /etc/mqtt-io.conf

