# hassio-addon-jethub-mqtt-io/rootfs/etc/services.d/mqtt-io/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting mqtt-io.."

# we need to use "exec" for correct stopping
# https://github.com/just-containers/s6-overlay/issues/229

cd /home/mqtt_io
exec /home/mqtt_io/venv/bin/python -m mqtt_io /etc/mqtt-io.conf

